# systemd/kingstonpier-cv-worker.service
[Unit]
Description=Kingston Pier CV worker (density counter -> readings)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__KP_USER__
WorkingDirectory=__KP_DIR__/tracker
# The worker APPENDS one count row per pass to this SQLite file; the API reads
# it. Set explicitly so it never depends on cwd-relative path resolution.
Environment=KP_DB_PATH=__KP_DIR__/db/data/readings.db
# Sample all five feeds every 3 minutes with the trained density model.
# Nice=10 so the inference bursts yield CPU to the rest of the Pi.
Nice=10
ExecStart=__KP_DIR__/tracker/.venv/bin/python crowd_tracker.py --watch 3
Restart=on-failure
RestartSec=10
# Light hardening (the worker only reads feeds + writes the readings DB).
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
